# roles/tomcat/files/tomcat9.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='WAPPtomcat:tomcat'>
<service name='application/tomcat' type='service' version='1'>
        <dependency name='filesystem-local'
                type='service'
                grouping='require_all'
                restart_on='none'>
                <service_fmri value='svc:/system/filesystem/local:default'/>
        </dependency>
        <dependency name='network'
                type='service'
                grouping='require_all'
                restart_on='refresh'>
                <service_fmri value='svc:/milestone/network:default'/>
        </dependency>
        <dependency name='name-services'
                type='service'
                grouping='require_all'
                restart_on='refresh'>
                <service_fmri value='svc:/milestone/name-services' />
        </dependency>
        <method_context>
                <method_credential user='tomcat' group='tomcat' />
        </method_context>
        <exec_method type='method' name='start'
                exec='/opt/tomcat/bin/startup.sh'
                timeout_seconds='300'>
        </exec_method> 
        <exec_method type='method' name='stop'
                exec='/opt/tomcat/bin/shutdown.sh stop'
                timeout_seconds='300' />
        <instance name='tomcat9' enabled='false'>
                <property_group name='conf' type='application'>
                        <propval name='catalina_home' type='astring' value='/opt/tomcat' />
                        <propval name='catalina_base' type='astring' value='/opt/tomcat' />
                </property_group>
        </instance>
        <stability value='Evolving' />
        <template>
                <common_name>
                        <loctext xml:lang='C'>
                                Tomcat application server
                        </loctext>
                </common_name>
        </template>
</service>
</service_bundle>
